# systemd/halo-agent.service
[Unit]
Description=halo-ai Agent — Autonomous Service Guardian
After=network.target
Wants=halo-llama-server.service halo-lemonade.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /srv/ai/agent/halo-agent.py
Restart=always
RestartSec=10
WatchdogSec=120

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=halo-agent
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
